# com.nyhetsagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nyhetsagent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/sigurd/ClaudeAgent/run.sh</string>
  </array>

  <!-- Kjører kl. 10:00 og 14:00 hver dag (lokal tid / Oslo) -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>10</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>14</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/sigurd/ClaudeAgent/logs/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/sigurd/ClaudeAgent/logs/launchd.log</string>
</dict>
</plist>
